# local-moonie.service
[Unit]
Description=Local Moonie OpenAI-compatible proxy with Ollama + Firefox search
After=network.target

[Service]
Type=simple
User=gabe
WorkingDirectory=/home/gabe/tmpcode/local_moonie
Environment="PATH=/usr/bin:/usr/local/bin"
EnvironmentFile=/home/gabe/tmpcode/local_moonie/.env
ExecStart=/usr/bin/python -m local_moonie.app.main
Restart=on-failure
RestartSec=5s

[Install]
WantedBy=multi-user.target

